# build/ansible-compute/roles/queue_processors/templates/queue_processors.service.j2
[Unit]
Description=Autoreduction queue processors service
After=network.target

[Service]
Type=forking
ExecStart=/home/{{ autoreduction_user }}/autoreduction/QueueProcessors/restart.sh
ExecReload=/home/{{ autoreduction_user }}/autoreduction/QueueProcessors/restart.sh
User={{ autoreduction_user }}
Group={{ autoreduction_user }}
Restart=always
RestartSec=9
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=activemq

[Install]
WantedBy=multi-user.target
