# build/ansible/roles/queue_processors/files/activemq.service
[Unit]
Description=ActiveMQ service
After=network.target

[Service]
Type=forking
ExecStart=/opt/autoreduce_deps/activemq/apache-activemq-5.15.5/bin/activemq start
ExecStop=/opt/autoreduce_deps/activemq/apache-activemq-5.15.5/bin/activemq stop
User=root
Group=root
Restart=always
RestartSec=9
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=activemq

[Install]
WantedBy=multi-user.target
